# etc/systemd/bluealsa-aplay.service.d/01-route-to-loopback.conf

[Service]
ExecStart=
ExecStart=/usr/bin/bluealsa-aplay -D hw:Loopback,0,0