# etc/systemd/zynthian.service
[Unit]
Description=Zynthian Sound Platform
After=jack2.service
Requires=jack2.service mod-ttymidi.service

[Service]
Environment=ZYNTHIAN_DIR=/zynthian
Environment=ZYNTHIAN_LOG_LEVEL=30
Environment=FRAMEBUFFER=/dev/fb1
Environment=HOME=/root
WorkingDirectory=/zynthian/zynthian-ui
ExecStart=/usr/bin/startx ./zynthian.sh -- :0 vt3
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
